# recipes-services/systemd/files/servicea.service
[Unit]
Description=SERVICE3

[Service]
ExecStart=/usr/bin/service_3

[Install]
WantedBy=multi-user.target
