# recipes-services/systemd/files/servicep.service
[Unit]
Description=SERVICE1

[Service]
ExecStart=/usr/bin/service_1

[Install]
WantedBy=multi-user.target
