# recipes-services/systemd/files/servicec.service
[Unit]
Description=SERVICE2

[Service]
ExecStart=/usr/bin/service_2

[Install]
WantedBy=multi-user.target
